# server/milliways-api.service
[Unit]
Description=Milliways API
After=network.target

[Service]
User=www-data
WorkingDirectory=/home/deploy/milliways/backend
Environment=DATABASE_URL=sqlite:////home/deploy/milliways/backend/milliways.db
Environment=FRONTEND_URL=https://milliways.nickgerrard.dev
ExecStart=/home/deploy/milliways/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8002
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
